# src/bin/patchmanager-daemon/systemd/checkForUpdates-org.SfietKonstantin.patchmanager.service
[Unit]
Description=Patchmanager check for updates service
Requires=dbus.service
After=dbus.service

[Service]
Type=oneshot
User=root
ExecStart=/usr/bin/env dbus-send --system --print-reply --dest=org.SfietKonstantin.patchmanager /org/SfietKonstantin/patchmanager org.SfietKonstantin.patchmanager.checkForUpdates
